# templates/systemd.service.erb
[Unit]
Description=<%= @service_name %> service
Documentation=https://github.com/thomasvincent/chef-cookbook-template
After=network.target

[Service]
Type=simple
User=<%= @user %>
Group=<%= @group %>
ExecStart=/usr/bin/<%= @service_name %> --port <%= @port %> --config /etc/example-service/service.conf
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log /var/lib/<%= @service_name %>
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target